# check_status.conf
[program:check_status]
command=php /var/www/html/Httpstatus/console.php controllers/internals/Httpstatus checkStatus
autostart=true
autorestart=true
startsecs=1
startretries=0
redirect_stderr=false