# supervisor.conf
[program:check_status]
command=php /var/www/html/httpstatus/console.php controllers/internals/Httpstatus checkStatus
autostart=true
autorestart=true
